# agents/main/skills/voice-trivia/trivia-voice.service
[Unit]
Description=Trivia Voice Server
After=network-online.target

[Service]
Type=simple
WorkingDirectory=%h/copylobsta/agents/main/skills/voice-trivia/server
# Edit WorkingDirectory if your repo is not at ~/copylobsta
ExecStart=/usr/bin/node index.js
Restart=always
RestartSec=5
EnvironmentFile=%h/.openclaw/.env

[Install]
WantedBy=default.target
