# agents/main/skills/copylobsta/copylobsta.service
[Unit]
Description=CopyLobsta Onboarding Server
After=network-online.target

[Service]
Type=simple
WorkingDirectory=%h/copylobsta/agents/main/skills/copylobsta/server
# Edit WorkingDirectory if your repo is not at ~/copylobsta
ExecStart=/usr/bin/node dist/index.js
Restart=always
RestartSec=5
EnvironmentFile=%h/.openclaw/.env

[Install]
WantedBy=default.target
